# crontab
* * * * * cd /var/www/html && php artisan queue:work --stop-when-empty && php artisan schedule:run >> /dev/stdout 2>&1